# etc/services.d/flask.ini
[program:flask]
directory=%(ENV_HOME)s/app
command=%(ENV_HOME)s/app/ENV/bin/uwsgi uwsgi.ini