# etc/services.d/nginx.ini
[program:nginx]
directory=%(ENV_HOME)s/nginx
command=%(ENV_HOME)s/bin/nginx -c %(ENV_HOME)s/nginx/conf/nginx.conf
